# scripts/stagedesk-power.service
[Unit]
Description=Stagedesk HDMI-CEC Power Control
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
ExecStart=/usr/bin/python3 /opt/stagedesk/stagedesk-power.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stagedesk-power

[Install]
WantedBy=multi-user.target
